# ansible/roles/dns-master/files/named.service
[Unit]
Description=Named

[Service]
Type=oneshot
ExecStart=/ansible/roles/dns-master/files/startup.script

[Install]
WantedBy=multi-user.target